# Linux/gfgapp.service
[Unit]
Description= My Python Flask App

[Service]
ExecStart= python3 /app.py
MemoryMax=100M
CPUShares=512